# .docker-php/8.4/rootfs/etc/service/horizon/run
#!/bin/sh
exec chpst -u app sh -lc 'cd /var/www/html && exec php artisan horizon'
